# emperor.uwsgi.service
[Unit]
Description=uwsgi emperor for shoppinglist website
After=network.target

[Service]
User=www-data
Restart=always
ExecStart=DIRECTORY/sl-env/bin/uwsgi --emperor DIRECTORY/sl-env/vassals --uid www-data --gid www-data

[Install]
WantedBy=multi-user.target
